# launchd/com.adamclarkson.content-vault.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.adamclarkson.content-vault</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>/Users/adamclarkson/dev/content-vault/scripts/run-local.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/adamclarkson/dev/content-vault</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
      <key>NODE_ENV</key>
      <string>production</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/adamclarkson/Library/Logs/content-vault.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/adamclarkson/Library/Logs/content-vault-error.log</string>
  </dict>
</plist>
